# configs/boxagent.service
[Unit]
Description=Simplified simple zebra service
After=syslog.target

[Service]
Type=simple
User=pi
Group=users
ExecStart=boxd server
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
